# infra/ops/systemd/bioconecta-backup.service
[Unit]
Description=BioConecta daily PostgreSQL backup
Requires=docker.service
After=docker.service

[Service]
Type=oneshot
ExecStart=/opt/bioconecta/infra/ops/bioconecta-backup.sh
